# crontab
MAILTO="derek.shaw@djshaw.ca"

# Grab an image every hour
0 * * * * OUTPUT=/home/cerf/djshaw.ca/diving/SYC /home/cerf/mygit/capture-syc/run.sh

# Retain the newest 7 image days based on filename day (YYYY-MM-DD), even
# when some calendar days are missing.
0 0 * * * /app/prune-old-images.sh /home/cerf/djshaw.ca/diving/SYC/ 7

